# projects/space-n-gators/launchd/com.tokengators.space-n-gators.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tokengators.space-n-gators</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Call webpack directly via node — avoids npm spawning an orphan child that launchd loses track of -->
    <string>/opt/homebrew/bin/node</string>
    <string>./node_modules/.bin/webpack</string>
    <string>serve</string>
    <string>--mode</string>
    <string>development</string>
    <string>--port</string>
    <string>8082</string>
    <string>--host</string>
    <string>0.0.0.0</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/operator/repos/MainFrame/projects/space-n-gators</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/operator/Library/Logs/space-n-gators.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/operator/Library/Logs/space-n-gators.error.log</string>
</dict>
</plist>
